# 99-test.rules
# rpi-midi-hub rule for test
ATTRS{idVendor}=="$2", ATTRS{idProduct}=="$3"
RUN+="/home/pi/bin/connectall.sh"
